# voice/voice.service
[Unit]
Wants=network.target
[Service]
ExecStart=/usr/local/bin/voice.sh
[Install]
WantedBy=default.target
